# provisioning/roles/base/templates/network-checker.service
[Unit]
Description=Log network status info

[Service]
Type=oneshot
ExecStart={{ script_name }}
